# root/etc/services.d/smartmontools/run
#!/usr/bin/with-contenv bash

exec /usr/sbin/smartd --debug
